# custom-packages/luci-app-iptv-manager/files/etc/init.d/iptv_manager
#!/bin/sh /etc/rc.common

START=95
STOP=15
USE_PROCD=1

PROG_MSD="/usr/bin/msd_lite"
CONFIG_NAME="iptv_manager"

_load_global() {
	config_get ENABLED global enabled  '0'
	config_get PROGRAM global program  'rtp2httpd'
}

# ── msd_lite：委托给其自身 init.d ────────────────────────────
# 原因：msd_lite 只接受 -c <config_file> 启动，不支持命令行参数配置。
#       其 init.d 负责读取 UCI，生成 JSON/XML 配置文件，再以 -c 启动二进制。

_start_msd() {
	if [ ! -f /etc/init.d/msd_lite ]; then
		logger -t iptv_manager -p daemon.err \
			"找不到 /etc/init.d/msd_lite，请检查 msd_lite 是否正确安装"
		return 1
	fi

	# 通过 UCI 启用 msd_lite（兼容不同 UCI section 命名方式）
	uci -q set msd_lite.@msd_lite[0].enabled='1' 2>/dev/null || \
	uci -q set msd_lite.msd_lite.enabled='1'     2>/dev/null || true
	uci -q commit msd_lite 2>/dev/null || true

	/etc/init.d/msd_lite start

	logger -t iptv_manager -p daemon.info \
		"启动 msd_lite（委托 /etc/init.d/msd_lite，配置由其自身读取）"
}

_stop_msd() {
	if [ -f /etc/init.d/msd_lite ]; then
		/etc/init.d/msd_lite stop
		# 重新禁用，防止 msd_lite 独立开机自启
		uci -q set msd_lite.@msd_lite[0].enabled='0' 2>/dev/null || \
		uci -q set msd_lite.msd_lite.enabled='0'     2>/dev/null || true
		uci -q commit msd_lite 2>/dev/null || true
	fi
	logger -t iptv_manager -p daemon.info "停止 msd_lite"
}

# ── rtp2HTTPd：委托给其自身 init.d ───────────────────────────
# 原因：rtp2httpd 二进制不直接读 UCI，由其 init.d 负责解析配置并启动。

_start_rtp2httpd() {
	if [ ! -f /etc/init.d/rtp2httpd ]; then
		logger -t iptv_manager -p daemon.err \
			"找不到 /etc/init.d/rtp2httpd，rtp2httpd feed 是否已安装？"
		return 1
	fi

	uci -q set rtp2httpd.@rtp2httpd[0].disabled='0'
	uci -q commit rtp2httpd

	/etc/init.d/rtp2httpd start

	logger -t iptv_manager -p daemon.info \
		"启动 rtp2HTTPd（委托 /etc/init.d/rtp2httpd，配置来自 /etc/config/rtp2httpd）"
}

_stop_rtp2httpd() {
	[ -f /etc/init.d/rtp2httpd ] && /etc/init.d/rtp2httpd stop

	uci -q set rtp2httpd.@rtp2httpd[0].disabled='1'
	uci -q commit rtp2httpd

	logger -t iptv_manager -p daemon.info "停止 rtp2HTTPd"
}

# ── procd 回调 ────────────────────────────────────────────────

start_service() {
	config_load "$CONFIG_NAME"
	_load_global

	[ "$ENABLED" = "1" ] || {
		logger -t iptv_manager -p daemon.info "IPTV Manager 未启用，跳过"
		return 0
	}

	case "$PROGRAM" in
		msd)       _start_msd       ;;
		rtp2httpd) _start_rtp2httpd ;;
		*)
			logger -t iptv_manager -p daemon.err "未知程序：$PROGRAM"
			;;
	esac
}

stop_service() {
	config_load "$CONFIG_NAME"
	_load_global

	case "$PROGRAM" in
		msd)       _stop_msd       ;;
		rtp2httpd) _stop_rtp2httpd ;;
	esac

	logger -t iptv_manager -p daemon.info "IPTV Manager 停止"
}

service_triggers() {
	procd_add_reload_trigger "$CONFIG_NAME"
	procd_add_reload_trigger "rtp2httpd"
}

reload_service() {
	stop
	start
}
