# custom-packages/luci-app-iptv-manager/files/etc/init.d/iptv_manager
#!/bin/sh /etc/rc.common

START=95
STOP=15
USE_PROCD=1

PROG_MSD="/usr/bin/msd_lite"
CONFIG_NAME="iptv_manager"

_load_global() {
	config_get ENABLED global enabled  '0'
	config_get PROGRAM global program  'rtp2httpd'
}

# ── rtp2HTTPd：委托给它自己的 init.d ─────────────────────────
# 原因：rtp2httpd 二进制不直接读 UCI，是其 init.d 负责
#       读取 /etc/config/rtp2httpd 并转成启动参数。
#       直接启动二进制会导致所有配置参数失效，使用默认值。

_start_rtp2httpd() {
	if [ ! -f /etc/init.d/rtp2httpd ]; then
		logger -t iptv_manager -p daemon.err \
			"找不到 /etc/init.d/rtp2httpd，rtp2httpd feed 是否已安装？"
		return 1
	fi

	# 将 rtp2httpd 自身 UCI 中的 disabled 改为 0
	# 这样它的 init.d 才会正常读配置并启动
	uci -q set rtp2httpd.@rtp2httpd[0].disabled='0'
	uci -q commit rtp2httpd

	/etc/init.d/rtp2httpd start

	logger -t iptv_manager -p daemon.info \
		"启动 rtp2HTTPd（委托 /etc/init.d/rtp2httpd，配置来自 /etc/config/rtp2httpd）"
}

_stop_rtp2httpd() {
	[ -f /etc/init.d/rtp2httpd ] && /etc/init.d/rtp2httpd stop

	# 重新将 disabled 设为 1，防止 rtp2httpd 独立开机自启（绕过我们的管理）
	uci -q set rtp2httpd.@rtp2httpd[0].disabled='1'
	uci -q commit rtp2httpd

	logger -t iptv_manager -p daemon.info "停止 rtp2HTTPd"
}

# ── msd_lite：procd 直接管理，参数来自 /etc/config/iptv_manager ──

_start_msd() {
	local interface port source_prefix buffer extra_args

	config_get interface     msd interface     'eth0'
	config_get port          msd port          '4022'
	config_get source_prefix msd source_prefix ''
	config_get buffer        msd buffer        '4096'
	config_get extra_args    msd extra_args    ''

	if [ ! -x "$PROG_MSD" ]; then
		logger -t iptv_manager -p daemon.err \
			"msd_lite 未安装：$PROG_MSD 不存在"
		return 1
	fi

	procd_open_instance "msd_lite"
	procd_set_param command "$PROG_MSD"
	procd_append_param command -i "$interface"
	procd_append_param command -p "$port"
	[ "$buffer" -gt 0 ] 2>/dev/null && \
		procd_append_param command -b "$buffer"
	[ -n "$source_prefix" ] && \
		procd_append_param command -s "$source_prefix"
	[ -n "$extra_args" ] && \
		procd_append_param command $extra_args
	procd_set_param respawn 3600 5 5
	procd_set_param stderr 1
	procd_close_instance

	logger -t iptv_manager -p daemon.info \
		"启动 msd_lite：接口=$interface 端口=$port"
}

# ── procd 回调 ────────────────────────────────────────────────

start_service() {
	config_load "$CONFIG_NAME"
	_load_global

	[ "$ENABLED" = "1" ] || {
		logger -t iptv_manager -p daemon.info "IPTV Manager 未启用，跳过"
		return 0
	}

	case "$PROGRAM" in
		msd)       _start_msd       ;;
		rtp2httpd) _start_rtp2httpd ;;
		*)
			logger -t iptv_manager -p daemon.err \
				"未知程序：$PROGRAM"
			;;
	esac
}

stop_service() {
	config_load "$CONFIG_NAME"
	_load_global

	# rtp2httpd 不由 procd 实例管理，需手动停止
	[ "$PROGRAM" = "rtp2httpd" ] && _stop_rtp2httpd

	# msd_lite 由 procd 自动停止其实例
	logger -t iptv_manager -p daemon.info "IPTV Manager 停止"
}

service_triggers() {
	procd_add_reload_trigger "$CONFIG_NAME"
	procd_add_reload_trigger "rtp2httpd"
}

reload_service() {
	stop
	start
}
